# config/redis.ini
[program:redis]
command=/usr/sbin/redis-server /etc/redis.conf
