# config/redis-sentinel.ini
[program:redis-sentinel]
command=/usr/sbin/redis-server /etc/sentinel.conf --sentinel
