# config/redis-test.ini
[program:redis-test]
command=/opt/bin/redis_test
